# focusflowd.in
#!/sbin/openrc-run

depend() {
    need dbus
}

description="FocusFlow user daemon"

export DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"
command_background="yes"

respawn_delay=5
respawn_max=0
pidfile="${XDG_RUNTIME_DIR:-/var/run}/focusflowd.pid"
supervisor="supervise-daemon"
command="@CMAKE_INSTALL_PREFIX@/bin/focusflowd"
